# etc/diod.service.in
[Unit]
Description=9P File Server

[Service]
Type=exec
ExecStart=@X_SBINDIR@/diod
ExecReload=kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
